# supervisord.conf
[supervisord]
nodaemon=true

[program:python]
command=/bin/bash -c "python /opt/aceproxy/acehttp.py"